# periodic/daily/backup
#!/bin/ash

mkdir -p /mnt/backup/daily
rsync –r --delete /mnt/source/ /mnt/backup/daily/
