# periodic/weekly/backup
#!/bin/ash

set -e # Exit on any error

. /usr/local/lib/backup-utils.sh

TIMESTAMP="date +\"%b %e %T %Z\""
BACKUP_DIR="${BACKUP:-/mnt/backup}"

if ! acquire_lock "weekly_backup"; then
    echo "$(eval $TIMESTAMP) ERROR: Another backup process is running" >>/var/log/backup
    exit 1
fi

# Ensure cleanup on exit
trap 'release_lock "weekly_backup"' EXIT

echo "$(eval $TIMESTAMP) Starting unison weekly backup..." >>/var/log/backup

# Ensure backup directory exists
if ! mkdir -p "$BACKUP_DIR/weekly"; then
    echo "$(eval $TIMESTAMP) ERROR: Failed to create weekly backup directory" >>/var/log/backup
    exit 1
fi

# Run unison and capture its exit status
if ! unison nextcloud_weekly >>/var/log/backup 2>&1; then
    echo "$(eval $TIMESTAMP) ERROR: Unison synchronization failed" >>/var/log/backup
    exit 1
fi

echo "$(eval $TIMESTAMP)  ... unison weekly backup finished" >>/var/log/backup
