# periodic/weekly/logrotate
#!/bin/ash

# Run logrotate
logrotate /etc/logrotate.conf
